# ansible/playbooks/roles/deploy/files/hpnssh.service.d/override.conf
[Unit]
# nginx is what signals "tamborg server ready"
Before=nginx.service
Wants=nginx.service
