# ansible/playbooks/roles/deploy/files/borg-daily.service
[Unit]
Description=Borg daily routine
Wants=borg-daily.timer

[Service]
Type=exec
User=borg
Group=borg
ExecSearchPath=/home/borg/.local/bin/
WorkingDirectory=/home/borg/tamborgcont/
ExecStart=poetry -q run python -m daily.main

[Install]
WantedBy=multi-user.target
